# skills/obs-dispatch/install/com.user.obs-dispatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.obs-dispatch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/opt/homebrew/bin/fswatch -0 --event Updated --event Created --event Renamed -e '\.obsidian/' -e '\.trash/' -e '\.git/' '/Users/zhiquntang/Documents/龙山军团' | xargs -0 -n1 /Users/zhiquntang/.local/bin/obs-dispatch.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/zhiquntang/Library/Logs/obs-dispatch.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zhiquntang/Library/Logs/obs-dispatch.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
